# deploy/myprobe.service
[Unit]
Description=MyProbe monitoring server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=myprobe
Environment=MYPROBE_LISTEN=0.0.0.0:25775
Environment=MYPROBE_DATABASE=/var/lib/myprobe/myprobe.db
EnvironmentFile=/etc/myprobe/server.env
ExecStart=/usr/local/bin/myprobe-server
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
